# files/etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command = /usr/sbin/php-fpm7.1
autorestart = true
stdout_logfile = /data/logs/php-fpm.log
stderr_logfile = /data/logs/php-fpm.log
